# overlay/etc/sv/ntpd/run
#!/bin/sh

set -eu

exec 2>&1
exec /usr/sbin/ntpd -n -l -p uk.pool.ntp.org
