# overlay/etc/sv/adsl/finish
#!/bin/sh

set -eu

exec 2>&1

pkill pppd			|| true
grep -q nas0: /proc/net/dev	&& ifdown nas0
