# overlay/etc/sv/syslogd/run
#!/bin/sh

set -eu

exec /sbin/syslogd -m 0
